# scripts/telegraf.upstart
# telegraf - The plugin-driven server agent for reporting metrics into InfluxDB
#
# See https://github.com/influxdata/telegraf

description     "Telegraf daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

setuid telegraf
script
    . /etc/default/telegraf
    exec /usr/bin/telegraf -config /etc/telegraf/telegraf.conf -config-directory /etc/telegraf/telegraf.d ${TELEGRAF_OPTS}
end script
