# plugins/inputs/supervisor/testdata/supervisord.conf
[supervisord]
nodaemon = true

[inet_http_server]
port = *:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:process-1]
command=/bin/sh -c 'while true; do echo "foo"; sleep 5; done'

[program:process-2]
command=/bin/sh -c 'while true; do echo "bar"; sleep 5; done'
