# plugins/inputs/supervisor/testdata/supervisord.conf
[inet_http_server]
port = 0.0.0.0:9001
username = login
password = pass

[supervisord]
logfile=/dev/stdout
pidfile=/var/run/supervisord.pid
logfile_maxbytes = 0
nodaemon = true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9001

[program:sleep]
process_name = %(program_name)s-%(process_num)s
command=/bin/sleep infinity
numprocs=3
autorestart=true